# node.conf
description "node.js server"
author      "joeyi"

# Start when disks mounted and network is up
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown
 
# Automatically Respawn:
respawn
respawn limit 99 5
 
script
    export HOME="/home/nodejs"
    exec sudo -u nodejs /usr/bin/node /home/nodejs/node_load.js >> /var/log/node.log 2>&1
end script
